# client/deploy/data/macos/AmneziaVPN.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>AmneziaVPN-service</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/AmneziaVPN.app/Contents/MacOS/AmneziaVPN-service</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>Sockets</key>
	<dict>
	    <key>Listeners</key>
	    <dict>
	        <key>SockServiceName</key>
                <string>5959</string>
	        <key>SockType</key>
	        <string>stream</string>
	        <key>SockFamily</key>
	        <string>IPv4</string>
	    </dict>
	</dict>
</dict>
</plist>
